# systemd/openclaw-restore.service
[Unit]
Description=Restore openclaw state from Dropbox
After=network-online.target
Wants=network-online.target
Before=openclaw.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/restore-openclaw-state.sh
TimeoutStartSec=60
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
