# systemd/post-boot-check.service
[Unit]
Description=Post-boot health check and config upgrade
After=openclaw.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/post-boot-check.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
